# deploy/adforge-mcp-web.service.example
[Unit]
Description=AdForge MCP web dashboard and API
After=network.target

[Service]
Type=simple
User=adforge
Group=adforge
WorkingDirectory=/opt/adforge-mcp
EnvironmentFile=/etc/adforge-mcp/adforge-mcp.env
Environment=AD_MCP_WEB_HOST=127.0.0.1
Environment=AD_MCP_WEB_PORT=8765
ExecStart=/opt/adforge-mcp/.venv/bin/ad-mcp-web
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/adforge-mcp /var/log/adforge-mcp

[Install]
WantedBy=multi-user.target
