# deploy/adforge-mcp-http.service.example
[Unit]
Description=AdForge MCP hosted HTTP transport
After=network.target adforge-mcp-web.service

[Service]
Type=simple
User=adforge
Group=adforge
WorkingDirectory=/opt/adforge-mcp
EnvironmentFile=/etc/adforge-mcp/adforge-mcp.env
Environment=AD_MCP_MCP_HTTP_HOST=127.0.0.1
Environment=AD_MCP_MCP_HTTP_PORT=8766
ExecStart=/opt/adforge-mcp/.venv/bin/ad-mcp-http
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/adforge-mcp /var/log/adforge-mcp

[Install]
WantedBy=multi-user.target
